# etc/mydb_backup.crontab
# crontab for mydb backup script goes here
MAILTO=""
07 18 * * * root docker exec -t dbaas bash -c 'DOCKER_CLIENT_TIMEOUT=3600; /opt/dbaas/backup_all.py' >/dev/null 2>&1 
